# launchagent/com.shift-ai.proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  SHIFT Proxy — macOS LaunchAgent

  Keeps the SHIFT preflight proxy running as a background service.
  All AI coding agents (OpenCode, Claude Code, Codex CLI, etc.) share
  this single proxy instance on port 8787.

  Install automatically:
    shift-ai setup

  Install manually:
    cp com.shift-ai.proxy.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.shift-ai.proxy.plist

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.shift-ai.proxy.plist
    rm ~/Library/LaunchAgents/com.shift-ai.proxy.plist

  Note: Update the ProgramArguments path to match your shift-ai installation.
  Default Homebrew path: /opt/homebrew/bin/shift-ai
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.shift-ai.proxy</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/shift-ai</string>
        <string>proxy</string>
        <string>start</string>
        <string>--port</string>
        <string>8787</string>
        <string>--foreground</string>
    </array>

    <!-- Restart if the proxy crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Start on login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/tmp/shift-proxy.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/shift-proxy.log</string>

    <!-- Environment -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
